# service-port-monitor.service
[Unit]
Description=Service Port Monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
Group=nogroup
WorkingDirectory=/opt/service-port-monitor
ExecStart=/opt/service-port-monitor/service-port-monitor
Restart=always
RestartSec=10

# Environment variables (alternatively use EnvironmentFile)
Environment="TELEGRAM_BOT_TOKEN=your_bot_token_here"
Environment="TELEGRAM_CHAT_ID=your_chat_id_here"
Environment="CONFIG_FILE=/opt/service-port-monitor/targets.txt"
Environment="LOG_FILE=/var/log/service-port-monitor/checker.log"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/service-port-monitor

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-port-monitor

[Install]
WantedBy=multi-user.target
